# files/systemd/pxls-backend.socket
[Socket]
ListenFIFO=%t/pxls/backend.stdin
Service=pxls-backend.service