# files/systemd/pxls-server.socket
[Socket]
ListenFIFO=%t/pxls.stdin
Service=pxls-server.service