# files/systemd/pxls-frontend.socket
[Socket]
ListenFIFO=%t/pxls/frontend.stdin
Service=pxls-frontend.service